# deadlock.se/websockthat/supervisor.achtung.kudzu.conf
[program:achtung]
priority=10
directory=/home/anton/deadlock.se/websockthat/
command=nodejs runserver.js achtung/config.js 2>&1 | tee -a achtung_server.log
user=root
autostart=true
autorestart=true
