# deadlock.se/websockthat/supervisor.pong.kudzu.conf
[program:pong]
priority=10
directory=/home/anton/deadlock.se/websockthat/
command=nodejs runserver.js pong/config.js 2>&1 | tee -a pong_server.log
user=root
autostart=true
autorestart=true
